# supervisord.conf
[supervisord]
logfile=lightning_logs/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10 ; (num of main logfile rotation backups;default 10)
loglevel=info ; (log level;default info; others: debug,warn,trace)
pidfile=lightning_logs/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false ; (start in foreground if true;default false)
minfds=1024 ; (min. avail startup file descriptors;default 1024)
minprocs=200 ; (min. avail process descriptors;default 200)
# user=root ; (default is current user, required if root)
# directory=/path/to/your/project ; (directory to cwd to before exec)

[program:exp]
command=python experiments/main.py
autostart=true
autorestart=true
stderr_logfile=lightning_logs/supervisord.err.log
stdout_logfile=lightning_logs/supervisord.out.log